# core/baselayout/files/rc.local
#!/bin/sh

# add your extra boot configs in here
